# argocd_sync/systemd/argocd-sheets-sync.service
[Unit]
Description=ArgoCD → Google Sheets sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/gongoeloe/projects/Ops_to_Biz
EnvironmentFile=/home/gongoeloe/projects/Ops_to_Biz/.env
ExecStart=/home/gongoeloe/projects/Ops_to_Biz/sync.sh
StandardOutput=journal
StandardError=journal
